# docker/var/spool/cron/crontabs/root
* * * * * /usr/local/bin/php /var/www/quiz2_cloudcomputing/artisan schedule:run
